# assets/etc/systemd/system/sys-kernel-config.mount
[Unit]
Description=Mount configfs

[Mount]
What=configfs
Where=/sys/kernel/config
Type=configfs

[Install]
WantedBy=multi-user.target
