# assets/yured.service
[Unit]
Description=yured
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/path/to/yured
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
